# com.example.twins-notification.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0//EN">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.example.twins-notification</string>

	<key>ProgramArguments</key>
	<array>
		<string>/PATH/TO/YOUR/twins-notification/twins-notification.zsh</string>
	</array>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>18</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<key>StandardOutPath</key>
	<string>/tmp/twins-launchd.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/twins-launchd.log</string>

	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
